# rootfs/etc/services.d/haproxy-up/run
#!/usr/bin/execlineb -P

with-contenv
exec haproxy-up
