# rootfs/etc/services.d/transform/run
#!/usr/bin/execlineb -P

with-contenv
consul-template -config /etc/consul.d/consul-template.cfg
